# rootfs/etc/s6-overlay/s6-rc.d/00-fix-permissions/run
#!/command/with-contenv sh
. "/usr/local/bin/logger"

echo "Setting permissions of /log..."
chown -R nobody:nobody /log
chmod -R 755 /log
