# rootfs/etc/s6-overlay/s6-rc.d/1000-wireguard-log/run
#!/bin/bash
program_name="wireguard"

s6-format-filter "%1 %s" "[${program_name}]" | logutil-service /log/"${program_name}"
