# rootfs/etc/s6-overlay/s6-rc.d/1000-wireguard/finish
#!/command/with-contenv bash

/bin/wg-quick down /run/wireguard/wg0.conf

/run/s6/basedir/bin/halt
