# rootfs/etc/s6-overlay/s6-rc.d/1000-wireguard/finish
#!/command/with-contenv bash

/usr/bin/wg-quick down wg0

/run/s6/basedir/bin/halt
